# doppelganger/deploy/com.nilsark.doppelganger.plist
<?xml version="1.0" encoding="UTF-8"?>
<!-- macOS launchd counterpart of doppelganger.service. Installed by install-macos.sh, which fills
     the {{...}} placeholders. Runs deploy/start.sh (self-update apply, then exec node). KeepAlive
     restarts the runtime on exit (the analog of systemd Restart=always). -->
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.nilsark.doppelganger</string>
  <key>ProgramArguments</key>
  <array>
    <string>{{REPO_DIR}}/doppelganger/deploy/start.sh</string>
  </array>
  <key>WorkingDirectory</key>
  <string>{{REPO_DIR}}/doppelganger</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>DOPPELGANGER_HOME</key>
    <string>{{DOPPELGANGER_HOME}}</string>
    <key>NODE_BIN</key>
    <string>{{NODE_BIN}}</string>
    <!-- launchd starts with a minimal PATH; include node's dir so `node`/`npm`/`git` resolve. -->
    <key>PATH</key>
    <string>{{NODE_DIR}}:/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>ThrottleInterval</key>
  <integer>5</integer>
  <key>StandardOutPath</key>
  <string>{{DOPPELGANGER_HOME}}/doppelganger.log</string>
  <key>StandardErrorPath</key>
  <string>{{DOPPELGANGER_HOME}}/doppelganger.log</string>
</dict>
</plist>
