# doppelganger/deploy/doppelganger.service
# Doppelgänger runtime — ONE process: dispatcher + scheduler. No OS cron.
# Installed as a systemd *user* unit by install.sh (claude/gws auth lives in the user's home).
# Runtime config (channels, ports, intervals) lives in $DOPPELGANGER_HOME/config.json; env vars
# below override it. install.sh scaffolds that file from config.example.json on first install.
[Unit]
Description=Doppelganger runtime (dispatcher + scheduler)
After=network-online.target

[Service]
Type=simple
WorkingDirectory={{REPO_DIR}}/doppelganger
# Self-update apply step (no-op unless selfUpdateEnabled). The leading "-" makes failure tolerable:
# an offline/failed update must never wedge startup — we just start the current code.
ExecStartPre=-{{REPO_DIR}}/doppelganger/deploy/update.sh
ExecStart={{NODE_BIN}} src/index.ts
Environment=DOPPELGANGER_HOME={{DOPPELGANGER_HOME}}
Restart=always
RestartSec=5

[Install]
WantedBy=default.target
